# LPIC2_02/freetime.service
[Unit]
Description=Free Memory and Time Display

[Service]
Type=oneshot
ExecStart=/usr/bin/freetime.sh

[Install]
WantedBy=graphical.target
